# data/deviceinfo-refresh.service
[Unit]
Description=Refresh device info
after=dbus.service

[Service]
ExecStart=device-info refresh

[Install]
WantedBy=multi-user.target
