# data/deviceinfo-refresh.conf
description "Refresh device info"

task

start on started dbus

exec device-info refresh
